# ops/LaunchAgents/com.cloudflare.bni-webhook-tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.cloudflare.bni-webhook-tunnel</string>
    <key>Comment</key>
    <string>Cloudflared named tunnel for BNI-Masta Recall.ai webhook (<your-webhook-host> → localhost:18821)</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ProgramArguments</key>
    <array>
      <string>/opt/homebrew/bin/cloudflared</string>
      <string>tunnel</string>
      <string>--config</string>
      <string>~/.cloudflared/config-bni.yml</string>
      <string>run</string>
    </array>
    <key>StandardOutPath</key>
    <string>~/.cloudflared/bni-tunnel.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>~/.cloudflared/bni-tunnel.stderr.log</string>
  </dict>
</plist>
